# gitfuse.service
[Unit]
Description=Mount git branches

[Service]
ExecStart=gitfuse

[Install]
WantedBy=default.target
